# supervisord.conf
[unix_http_server]
file=/Users/john/workspace/tmp/django-datatable/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)

[supervisord]
logfile=/Users/john/workspace/tmp/django-datatable/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/Users/john/workspace/tmp/django-datatable/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
umask=022                   ; (process file creation umask;default 022)
identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
directory=/Users/john/workspace/python/django-datatable             ; (default is not to cd during start)
nocleanup=false              ; (dont clean up tempfiles at start;default false)
childlogdir=/Users/john/workspace/tmp/django-datatable/tmp            ; ('AUTO' child log dir, default $TEMP)
strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/john/workspace/tmp/django-datatable/tmp/supervisor.sock ; use a unix:// URL  for a unix socket
prompt=supervisor         ; cmd line prompt (default "supervisor")
history_file=/Users/john/workspace/tmp/django-datatable/.sc_history  ; use readline history if available

[program:django_datatable]
command=/Users/john/.virtualenvs/datatable/bin/gunicorn django_datatable_demo.wsgi:application
        --workers 4
        --max-requests 10000
        -k gevent
        --timeout 30
        --log-level info
        --bind 0.0.0.0:1115%(process_num)1d

process_name=%(program_name)s_%(process_num)d
umask=022
startsecs=10
stopwaitsecs=0
redirect_stderr=true
stdout_logfile=/Users/john/workspace/tmp/django-datatable/logs/process_%(process_num)02d.log
numprocs=4
numprocs_start=0
